# dynamiclock.timer
[Unit]
Description=Dynamic Lock timer.

[Timer]
OnActiveSec=5s
# 每次运行后再次等待5秒
OnUnitActiveSec=5s  
Persistent=true

[Install]
WantedBy=timers.target